# etc/profile.d/lang.csh.new
#!/bin/csh

setenv LANG ru_RU.utf8
setenv LC_NUMERIC C
setenv LC_COLLATE C
