# etc/profile.d/lang.sh.new
#!/bin/sh

export LANG=ru_RU.utf8
export LC_NUMERIC=C
export LC_COLLATE=C
